# .devcontainer/supervisord.conf
[supervisord]
nodaemon=false
logfile=/workspace/logs/supervisord.log
pidfile=/workspace/logs/supervisord.pid
childlogdir=/workspace/logs

[program:django]
command=gunicorn job_applier.wsgi:application --bind 0.0.0.0:8000 --workers 2 --threads 4 --worker-class gthread --timeout 120
directory=/workspace
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/workspace/logs/django.log
environment=PATH="/usr/local/bin:/usr/bin:/bin",DJANGO_SETTINGS_MODULE="job_applier.settings"
user=vscode

[program:celery-worker]
command=celery -A job_applier worker -l info --concurrency=2
directory=/workspace
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/workspace/logs/celery-worker.log
environment=PATH="/usr/local/bin:/usr/bin:/bin",DJANGO_SETTINGS_MODULE="job_applier.settings"
user=vscode
stopwaitsecs=60
stopasgroup=true
killasgroup=true

[program:celery-beat]
command=celery -A job_applier beat -l info
directory=/workspace
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/workspace/logs/celery-beat.log
environment=PATH="/usr/local/bin:/usr/bin:/bin",DJANGO_SETTINGS_MODULE="job_applier.settings"
user=vscode

[program:flower]
command=celery -A job_applier flower --port=5555 --basic_auth=%(ENV_FLOWER_BASIC_AUTH)s
directory=/workspace
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/workspace/logs/flower.log
environment=PATH="/usr/local/bin:/usr/bin:/bin",DJANGO_SETTINGS_MODULE="job_applier.settings"
user=vscode

[group:all]
programs=django,celery-worker,celery-beat

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
chown=vscode:vscode

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock